# gunicorn_Doggy.socket
[Unit]
Description=Doggy socket
[Socket]
ListenStream= /home/wwwadmin/Doggy/Doggy.sock 
[Install]
WantedBy=sockets.target
